# kafka/supervisor/confluent-schema.conf
[program:confluentSchema]
command=/confluent/bin/schema-registry-start /confluent/etc/schema-registry/schema-registry.properties
autostart=true
autorestart=true
