# kafka/supervisor/kafka.conf
[program:kafka]
command=/confluent/bin/kafka-server-start /confluent/etc/kafka/server.properties
autostart=true
autorestart=true
