# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/confluent/bin/zookeeper-server-start /confluent/etc/kafka/zookeeper.properties
autostart=true
autorestart=true
